# LaunchAgent/com.sai.notchly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sai.notchly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/saikiran/openclaw/notchly_v2/build/NotchlyV2.app/Contents/MacOS/NotchlyV2</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/saikiran/notchly/v2/logs/notchly.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/saikiran/notchly/v2/logs/notchly_err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/saikiran</string>
    </dict>
</dict>
</plist>
