# LaunchAgent/com.sai.notchly.brain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sai.notchly.brain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/saikiran/openclaw/notchly_v2/Python/brain_loop.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/saikiran/notchly/v2/logs/brain.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/saikiran/notchly/v2/logs/brain_err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/saikiran</string>
        <key>PYTHONPATH</key>
        <string>/Users/saikiran/openclaw/notchly_v2/Python</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/saikiran/openclaw/notchly_v2/Python</string>
</dict>
</plist>
